# crontab
@reboot /home/setupini/akvarium.sh
* * * * * /usr/bin/php /var/www/html/crontab.php
